# deploy/roles/deploy/templates/mqtt2pg.service
[Unit]
Description=MQTT2PG Service
After=syslog.target
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=simple
Restart=always
RestartSec=60
WorkingDirectory={{workingdir}}
User={{user}}
Group={{user}}
Environment=MQTT_URL=mqtt://localhost
Environment=MQTT_USER={{user}}
Environment=MQTT_PASSWORD={{mqtt_password}}
Environment=PG_URL=postgresql://{{user}}:{{pg_password}}@localhost/hab
ExecStart={{workingdir}}/mqtt2pg run

[Install]
WantedBy=multi-user.target
